# root/etc/systemd/system/article-scripts.timer
[Unit]
Description=Run article scripts hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target